# docker/php-worker/supervisor/conf.d/websockets.conf
[program:websockets]
command=php /var/www/app/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=app
